# scripts/sierra.service
[Unit]
Description=Sierra smart irrigation stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__SIERRA_DIR__
User=__USER__
Group=docker
ExecStart=/usr/bin/docker compose -f __SIERRA_DIR__/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f __SIERRA_DIR__/docker-compose.yml down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
